# components/portal/deploy/raiccoon-portal.service
[Unit]
Description=Your Organization Internal Portal
After=network.target

[Service]
Type=simple
User=operator
WorkingDirectory=/opt/raiccoon/internal-secops-portal
Environment=SECOPS_PORTAL_HOST=0.0.0.0
Environment=SECOPS_PORTAL_PORT=8008
Environment=SECOPS_ALLOWED_CIDRS=127.0.0.1/32,10.0.0.0/24
Environment=PORTAL_SECRET_KEY=replace-with-long-random-secret
Environment=RAICCOON_ADMIN_USERNAME=admin
Environment=RAICCOON_ADMIN_PASSWORD=change-me-now
Environment=HERMES_CLI=/opt/raiccoon/.local/bin/hermes
ExecStart=/opt/raiccoon/internal-secops-portal/.venv/bin/python /opt/raiccoon/internal-secops-portal/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
